# scripts/systemd/netbox-agent-test-network.service
[Unit]
Description=NetBox Agent Test — Network-Only Sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/netbox-agent-test
Environment=PYTHONPATH=/opt/netbox-agent-test
ExecStart=/usr/bin/python3 -m netbox_agent.cli -c /etc/netbox-agent-test/config.yml --network-only
StateDirectory=netbox-agent-test
StandardOutput=journal
StandardError=journal
SyslogIdentifier=netbox-agent-test-network

[Install]
WantedBy=multi-user.target
