# scripts/systemd/netbox-agent-test-daily.timer
[Unit]
Description=NetBox Agent Test — Daily Full Sync Timer

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=300
Unit=netbox-agent-test.service

[Install]
WantedBy=timers.target
