# scripts/systemd/netbox-agent-test-boot.service
[Unit]
Description=NetBox Agent Test — Boot Sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/netbox-agent-test
Environment=PYTHONPATH=/opt/netbox-agent-test
ExecStartPre=/bin/sleep 60
ExecStart=/usr/bin/python3 -m netbox_agent.cli -c /etc/netbox-agent-test/config.yml --update-all
StateDirectory=netbox-agent-test
StandardOutput=journal
StandardError=journal
SyslogIdentifier=netbox-agent-test-boot

[Install]
WantedBy=multi-user.target
